# supervisor/conf.d/faucet.conf
[group:faucet]
programs=faucet

[program:faucet]
command=/root/.aen/bin/faucet-startup
autorestart=true
